# src/php/cli/7.3/rootfs/etc/periodic/daily/security-checker
#!/bin/sh

COMPOSER_LOCK_PATH="/app/composer.lock"
if test -f "$COMPOSER_LOCK_PATH"; then
  security-checker security:check ${COMPOSER_LOCK_PATH} 2>&-
fi
