# src/php/fpm/7.3/rootfs/etc/supervisor.d/php-fpm.ini
[group:php-fpm]
programs = php-fpmd
priority = 20

[program:php-fpmd]
command = /usr/sbin/php-fpm7 -F
process_name = %(program_name)s
numprocs = 1
user = root
startsecs = 0
autostart = true
autorestart = true
exitcodes = 0
stopsignal = TERM
restartsignal = HUP
reloadsignal = USR2
stdout_events_enabled = false
stdout_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
redirect_stderr = true
